# ver1_0/jobs/crontab
MAILTO="fragro@gmail.com"
59 23 * * * ~/current/jobs/rebuild_index
